# deploy/feeder_supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[supervisord]

[program:celery]
process_name = celery-%(process_num)s
command = /home/feeder/.virtualenvs/feeder/bin/python -m feedreader.tasks.worker --number=%(process_num)s --config=deploy/production.yaml
directory=/home/feeder/project
stdout_logfile=/home/feeder/worker.log
stderr_logfile=/home/feeder/worker.err
autorestart=true
numprocs = 1

[program:tornado]
process_name = feeder-%(process_num)s
command=/home/feeder/.virtualenvs/feeder/bin/python -m feedreader.main --updates --port 1000%(process_num)s --number=%(process_num)s --config=deploy/production.yaml
directory=/home/feeder/project
stdout_logfile=/home/feeder/server.log
stderr_logfile=/home/feeder/server.err
autorestart=true
numprocs = 1
